# services/archiver/matrix-archiver.service
[Unit]
Description=Matrix Message Archiver for Supabase
After=network.target
Wants=network.target

[Service]
Type=simple
User=matrix
Group=matrix
WorkingDirectory=/home/matrix-ai/services/archiver
ExecStart=/usr/bin/python3 -m src.archiver.main
Restart=always
RestartSec=10
Environment=PYTHONPATH=/home/matrix-ai/services/archiver

# Environment variables (set these with your actual values)
Environment=MATRIX_PASSWORD=your_archiver_password
Environment=SUPABASE_URL=your_supabase_url
Environment=SUPABASE_SERVICE_ROLE_KEY=your_service_role_key
Environment=SUPABASE_ANON_KEY=your_anon_key
Environment=SUPABASE_JWT_SECRET=your_jwt_secret
Environment=SUPABASE_DB_URL=your_database_url

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=matrix-archiver

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/matrix-ai/logs

[Install]
WantedBy=multi-user.target